# roles/application/templates/nodeapp.conf
start on filesystem and started networking
respawn
chdir /var/www/DemoProject  #deployment directory
exec /usr/bin/node app.js  #start command - no forever needed, if it fails, the service restarts
